# roles/automatic_turn_off/templates/auto-off.service.j2
# {{ ansible_managed }}

[Unit]
Description=Automatic Shutdown Service
After=network.target

[Service]
Type=oneshot
ExecStart=/bin/systemctl poweroff

[Install]
WantedBy=multi-user.target
